# templates/appdyn_machineagent.sh.erb
#!/bin/sh
#
# AppDynamics Machine Agent init script
#
# chkconfig: - 98 3
# description: AppDynamics Machine Agent

# Source function library.
. /etc/init.d/functions
. /etc/profile

RETVAL=0

JAVA="/usr/bin/java"

AGENT_HOME="/opt/appdynamics/MachineAgent-<%= @instance %>"
AGENT="$AGENT_HOME/machineagent.jar"

AGENT_OPTIONS="$AGENT_OPTIONS -Xmx10m"

AGENT_LOG="$AGENT_HOME/machineagent.log"

prog="AppDynamics Machine Agent"

start() {
        # Start daemon.
        echo -n $"Starting $prog"
        daemon --user=appdyn "$JAVA $AGENT_OPTIONS -jar $AGENT &"
        # nohup java -Dappdynamics.agent.uniqueHostId=source-backend.ac-jas02.acce.alfa.local -jar machineagent.jar &
        RETVAL=$?
        [ $RETVAL -eq 0 ]
        return $RETVAL
}

stop() {
        # Stop daemon.
        action  $"Shutting down $prog: " pkill -f machineagent.jar
        RETVAL=$?
        [ $RETVAL -eq 0 ]
        return $RETVAL
}

status() {
        # Daemon status.
        pgrep -f machineagent.jar 
        RETVAL=$?
        return $RETVAL
}


# See how we were called.
case "$1" in
  start)
        start
        ;;
  stop)
        stop
        ;;
  status)
        status
        ;;
  restart|reload)
        stop
        sleep 5s
        start
        RETVAL=$?
        ;;
  *)
        echo $"Usage: $0 {start|stop|restart|status}"
        exit 1
esac

exit $RETVAL
